# deployment/systemd/malaika_nest-health.timer
[Unit]
Description=Run Malaika Nest healthcheck every 15 minutes

[Timer]
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
