# deployment/systemd/malaika_nest-backup.timer
[Unit]
Description=Run Malaika Nest backup daily

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target