# deployment/systemd/malaika_nest-log-monitor.service
[Unit]
Description=Malaika Nest Auth Failure Monitor
After=network.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/malaikanest/monitor.env
ExecStart=/bin/bash /home/mohameddosho20/malaikanest/deployment/log-monitor.sh