# deployment/systemd/malaika_nest-renew.service
[Unit]
Description=Malaika Nest cert renewal runner
After=docker.service

[Service]
Type=oneshot
User=ubuntu
# Set REPO_DIR when enabling the unit or via /etc/default/malaika_nest
Environment=REPO_DIR=/home/ubuntu/malaika_nest
WorkingDirectory=${REPO_DIR}
ExecStart=/bin/bash ${REPO_DIR}/deployment/renew_certs.sh
Nice=10
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
