# deployment/celerybeat.service
[Unit]
Description=Celery Beat Scheduler for Malaika Nest
Documentation=https://docs.celeryq.dev
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/malaikanest/backend
EnvironmentFile=-/var/www/malaikanest/backend/.env.production
Environment=DJANGO_ENV=prod
Environment=DJANGO_SETTINGS_MODULE=config.settings.prod

ExecStart=/var/www/malaikanest/venv/bin/celery \
    -A config \
    beat \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler \
    --loglevel=info \
    --logfile=/var/log/malaikanest/celery-beat.log \
    --pidfile=/run/celery/celerybeat.pid

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=celerybeat-malaikanest

[Install]
WantedBy=multi-user.target
